# s6/nginx/run
#!/usr/bin/env bash
#
# Jenkins s6 run script
#
NGINX=$(which nginx)
echo 'Starting nginx from ' $NGINX
$NGINX
# -g 'pid /run/nginx.pid; error_log stderr;'"